# systemd/multiagentedustack-synthesis.timer
[Unit]
Description=DEPRECATED: use multiagentedustack-digest.timer (daily) and multiagentedustack-weekly.timer (Sunday)

[Timer]
# Historically Sunday-only; now daily so legacy enables still get article digests.
# Disable this unit after enabling multiagentedustack-digest.timer.
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
