# systemd/multiagentedustack-weekly.timer
[Unit]
Description=Run MultiAgentEDUstack weekly wiki + forecast

[Timer]
# Sunday 03:30 -- after Saturday's daily digest, ahead of the work week.
OnCalendar=Sun *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
